# systemd/systemd-echo.socket.in
[Unit]
Description=ECHO Socket

[Socket]
ListenStream=7
ListenDatagram=7

[Install]
WantedBy=sockets.target
